# files/etc/hotplug.d/iface/50-ap-isolation
#!/bin/sh

[ "$ACTION" = "ifup" ] || exit 0

case "$DEVICE" in
	phy*|wlan*) /etc/init.d/ap-isolation reload >/dev/null 2>&1 & ;;
	*) exit 0 ;;
esac
